# supervisord/conf/supervisor/worker.conf
[program:worker]
command=docker exec darbit_docker_php php /var/www/bitrix/local/modules/phmk_corezoid/scripts/supervisord.php
numprocs=1
autostart=true
autorestart=true
user=root