# php/php74/conf/supervisor/worker.conf
[program:worker]
command=/usr/bin/php /var/www/bitrix/local/modules/phmk_corezoid/scripts/supervisord.php
numprocs=1
autostart=true
autorestart=true
user=root